# docs/launchd/com.grounded-weather-forecast.maintain.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <!-- The nightly retrain: rebuild -> backtest -> report -> truth QC.
       This is what makes "improves as the archive grows" happen unattended;
       promotions land in the append-only release ledger. -->
  <key>Label</key>
  <string>com.grounded-weather-forecast.maintain</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>-lc</string>
    <string>cd __REPO__ &amp;&amp; grounded-weather-forecast build-dataset &amp;&amp; grounded-weather-forecast backtest --source live &amp;&amp; grounded-weather-forecast report &amp;&amp; grounded-weather-forecast truth-qc</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>SYNOPTIC_TOKEN</key>
    <string>__SYNOPTIC_TOKEN__</string>
  </dict>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>2</integer>
    <key>Minute</key>
    <integer>15</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>__LOG_DIR__/maintain.log</string>
  <key>StandardErrorPath</key>
  <string>__LOG_DIR__/maintain.log</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
